# deploy/systemd/cloud-agents-worker.service
[Unit]
Description=aflow Remote Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agentflow
EnvironmentFile=/etc/cloud-agents-worker.env
Environment=PYTHONPATH=/opt/agentflow/runtime
ExecStart=/usr/bin/python3 -m cloud_agents_runtime.worker
Restart=always
RestartSec=5
User=cloudagents
Group=cloudagents
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/cloud-agents-worker
CPUAccounting=true
CPUQuota=100%
MemoryAccounting=true
MemoryMax=1G
TasksAccounting=true
TasksMax=512

[Install]
WantedBy=multi-user.target
